# config/includes.chroot/etc/cron.d/amfreos-autoenroll
@reboot root /usr/local/amfreos/auto-enroll.sh
